# config/tomcat.conf
[program:tomcat]
command=/opt/tomcat/bin/startup.sh
